# downloads-folder-cleaner/crontab
@daily python3 -u "/home/dan/Desktop/dev/py/python-mini-projects/downloads-cleaner/main.py"