# downloads-folder-cleaner/crontab
@daily python3 -u "/home/dan/Desktop/dev/py/downloads-cleaner/main.py"